# URL Classifier/deploy/ec2/cleannet-classifier.service
[Unit]
Description=CleanNet URL Classifier (Gunicorn)
After=network.target

[Service]
Type=simple
# Ubuntu EC2: use User=ubuntu Group=ubuntu. Amazon Linux: ec2-user.
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/url-classifier
Environment="PATH=/opt/url-classifier/venv/bin"
Environment="CLEANNET_DISCOVERY=0"
Environment="TRUST_PROXY=1"
Environment="FLASK_DEBUG=0"
Environment="GUNICORN_BIND=127.0.0.1:5000"
Environment="GUNICORN_WORKERS=4"
ExecStart=/opt/url-classifier/venv/bin/gunicorn -c gunicorn.conf.py wsgi:app
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
